# com.davidmerz.software-kitchen-ticketer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davidmerz.software-kitchen-ticketer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/davidmerz/projects/Software-Kitchen-Ticketer/software-kitchen-ticketer.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/davidmerz/projects/Software-Kitchen-Ticketer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/davidmerz/Library/Logs/software-kitchen-ticketer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davidmerz/Library/Logs/software-kitchen-ticketer.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
